# cron.daily/ribasis-archive
#!/bin/bash

cd /root/ribasis/rib-receiver
/root/ribasis/rib-receiver/venv/bin/python /root/ribasis/rib-receiver/archive.py
